# tools/sonarqube/supervisor/jobs/sonarqube.ini
[program:sonarqube]
command = bash -c "${SONARQUBE_HOME}/bin/linux-x86-64/sonar.sh start"
priority = 100
startsecs = 0
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/inovedados/stdout.log
stderr_logfile = /var/log/inovedados/stderr.log
stdout_events_enabled = true