# php/apache/supervisor/jobs/crond.ini
[program:crond]
command = crondtab
startsecs = 0
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/cron/stdout.log
stderr_logfile = /var/log/cron/stderr.log
stdout_events_enabled = true