# php/apache/supervisor/jobs/crond.ini
[program:crond]
command=crond -n
startsecs=3
stopwaitsecs=3
autostart=true
autorestart=true
stdout_logfile=/var/log/cron/stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/cron/stderr.log
stderr_logfile_maxbytes=10MB
user=inovedados
