# php/apache/supervisor/jobs/httpd.ini
[program:httpd]
command=exec httpd -DFOREGROUND
autostart=true 
autorestart=true 
stdout_logfile=/var/log/httpd/stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/httpd/stderr.log
stderr_logfile_maxbytes=10MB
user=inovedados
exitcodes=0